# android-app/desktopApp/packaging/localagent.user.service
# Local Agent — headless desktop agent (systemd USER service).
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp localagent.user.service ~/.config/systemd/user/localagent.service
#   # edit ExecStart below to match how you installed the app (see DESKTOP_PACKAGING.md)
#   systemctl --user daemon-reload
#   systemctl --user enable --now localagent.service
#   sudo loginctl enable-linger "$USER"      # so it runs at boot before you log in
#   journalctl --user -u localagent -f      # follow logs

[Unit]
Description=Local Agent (headless desktop agent)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=LOCALAGENT_HEADLESS=1
# Optional: force the CPU llama-server build instead of the Vulkan GPU one.
# Environment=LOCALAGENT_LLAMA_SERVER_VARIANT=cpu
#
# Secrets (Brave/HF keys, keystore password): prefer provisioning via the app's
# SecureStorage (run the GUI once) so the service needs NO secret env vars. If you must
# pass one, do NOT inline `Environment=BRAVE_API_KEY=...` — a process's env is readable via
# /proc/<pid>/environ. Use a 0600 file instead (see "Secrets & /proc exposure" in
# DESKTOP_PACKAGING.md):
#   EnvironmentFile=%h/.config/localagent/secrets.env   # 0600, KEY=VALUE lines
# .deb install path shown; for a `createDistributable` app-image use that bin/LocalAgent.
ExecStart=/opt/local-agent/bin/LocalAgent
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
